# data/ModemManager.service.in
[Unit]
Description=Modem Manager
After=syslog.target

[Service]
Type=dbus
Environment=LD_LIBRARY_PATH=/usr/local/lib
BusName=org.freedesktop.ModemManager1
ExecStart=@sbindir@/ModemManager --debug
StandardError=null
Restart=on-abort

[Install]
WantedBy=multi-user.target
Alias=dbus-org.freedesktop.ModemManager1.service
